# amazon-mcp.service
[Unit]
Description=Amazon SP API MCP Server
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/amazon-mcp-server
Environment=PATH=/opt/amazon-mcp-server/venv/bin
Environment=PYTHONPATH=/opt/amazon-mcp-server
ExecStart=/opt/amazon-mcp-server/venv/bin/python amazon_mcp_http_server.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=amazon-mcp

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/amazon-mcp-server

[Install]
WantedBy=multi-user.target